# deconz/rootfs/etc/services.d/deconz/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: AirSonos
# Runs the AirSonos server
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

declare -a options

hass.log.info 'Starting the deCONZ server'

if hass.config.true 'xpra'; then
  xpra start --bind-tcp=0.0.0.0:8236 :100
  export DISPLAY=:100
else
  options+=(-platform minimal)
fi

# Non-interactive

# Connect to zigbee adapter
options+=(--auto-connect=1)

options+=(--dbg-info=1)

# Bind to a specific interface
if hass.config.has_value 'httpport'; then
    options+=("--http-port=$(hass.config.get 'httpport')")
fi
if hass.config.has_value 'wsport'; then
    options+=("--ws-port=$(hass.config.get 'wsport')")
fi

# Run deCONZ
exec deCONZ "${options[@]}"
